# deploy/betting-helper.service
[Unit]
Description=Betting Helper FastAPI App
After=network.target

[Service]
User=root
WorkingDirectory=/opt/betting_helper
ExecStart=/opt/betting_helper/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8765
Restart=always

[Install]
WantedBy=multi-user.target
